# systemd/.config/systemd/user/low-battery-check.service
[Unit]
Description=Low Battery Notification

[Service]
ExecStart=%h/.local/bin/low-battery-check

[Install]
WantedBy=default.target
