# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=Udiskie
Documentation=man:udiskie(8)

[Service]
ExecStart=/usr/bin/udiskie
KillSignal=SIGINT

[Install]
WantedBy=default.target
