# docker/.config/systemd/user/docker-system-prune.service
[Unit]
Description=Prune unused docker data

[Service]
ExecStart=/usr/bin/env docker system prune --force

[Install]
WantedBy=default.target
